# config-management/roles/harrow.backend/files/etc/systemd/system/harrow-notify-about-failure@.service
[Unit]
Description=Notify on systemd failures
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/notify-about-failure "%I" "Hostname: %H" "Machine ID: %m" "Boot ID: %b"

