# config-management/roles/harrow.backend/templates/etc/systemd/system/harrow-api.service.j2
[Unit]
Description=Harrow API
After=harrow.service
{% if harrow.services.notify_on_failure %}
OnFailure=harrow-notify-about-failure@%n.service
{% endif %}

[Service]
EnvironmentFile=/etc/harrow/env
WorkingDirectory=/tmp
PrivateTmp=true
ExecStart=/usr/local/bin/harrow api
User=harrow
Restart=always
