# config-management/roles/harrow.backend/templates/etc/systemd/system/harrow-activity-worker.service.j2
[Unit]
Description=Harrow Activity Worker
#Requires=harrow.service
After=harrow.service
{% if harrow.services.notify_on_failure %}
OnFailure=harrow-notify-about-failure@%n.service
{% endif %}

[Service]
EnvironmentFile=/etc/harrow/env
WorkingDirectory=/tmp
PrivateTmp=true
ExecStart=/usr/local/bin/harrow activity-worker
User=harrow
Restart=always
RestartSec=5
