# config-management/roles/harrow.backend/files/etc/profile.d/harrow.sh
export GOPATH=/harrow
export PATH=$PATH:/harrow/bin
